# installer/desktop_installer/resources/zsnap-5min.service
[Unit]
Description=Takes a snapshot of home dataset

[Service]
Type=simple
ExecStart=/sbin/autosnap rpool/HOME 5min 20
